# rc.d/spkrd
#!/bin/sh

# FreeBSD rc.d script for spkrd - FreeBSD speaker device network server
#
# Add the following lines to /etc/rc.conf to enable spkrd:
#
# spkrd_enable="YES"
# spkrd_flags="--port 8080 --device /dev/speaker --retry-timeout 30"
#
# Available flags:
#   --port <port>           Server port (default: 8080)
#   --device <path>         Speaker device path (default: /dev/speaker)  
#   --retry-timeout <secs>  Device retry timeout (default: 30)
#   --daemon                Run as daemon (automatically added by rc.d)
#   --pidfile <path>        PID file path (default: /var/run/spkrd.pid)
#   --debug/-d              Enable debug logging including client requests
#
# Examples:
#   spkrd_flags="--port 3000"
#   spkrd_flags="--device /tmp/test-speaker --port 9000"
#   spkrd_flags="--retry-timeout 60 --port 8080"
#   spkrd_flags="--pidfile /var/run/spkrd-alt.pid --port 8080"
#   spkrd_flags="--debug"                    # Enable debug logging
#   spkrd_flags="--debug --port 8080"       # Debug logging with custom port
#
# Logging:
#   - Default: Startup messages and errors only (to syslog facility 'daemon')
#   - With --debug: Also logs client requests with IP, melody data, and retry counts
#   - View logs: tail -f /var/log/daemon.log | grep spkrd
#   - Non-root users: Use custom --pidfile (e.g., /tmp/spkrd.pid)

. /etc/rc.subr

name=spkrd
rcvar=spkrd_enable

command="/usr/local/bin/${name}"
command_args="--daemon --pidfile ${pidfile} ${spkrd_flags}"

pidfile="/var/run/${name}.pid"
required_files="/usr/local/bin/${name}"

load_rc_config $name
run_rc_command "$1"
