# systemd/spkrd.service
# systemd user unit for spkrd - speaker device network server
#
# Installed as a user service so spkrd has access to the per-user
# PulseAudio / PipeWire socket. Each user manages their own instance.
#
# Install to /usr/local/lib/systemd/user/ (or wherever DSTDIR resolves to):
#   make install
#
# Then enable for the current user:
#   systemctl --user daemon-reload
#   systemctl --user enable spkrd
#   systemctl --user start spkrd
#
# To start automatically on boot even when not logged in:
#   loginctl enable-linger $USER
#
# To override flags without editing this file, create a drop-in:
#   systemctl --user edit spkrd
# and add an [Service] section, e.g.:
#   [Service]
#   ExecStart=
#   ExecStart=/usr/local/bin/spkrd --port 3000 --debug
#
# Logs are captured by journald:
#   journalctl --user -u spkrd -f

[Unit]
Description=spkrd speaker device network server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/spkrd --port 1111 --cpal-host=PulseAudio
Restart=on-failure

[Install]
WantedBy=default.target
